# webui/setup/supervisord.conf
[supervisord]
http_port=/var/tmp/supervisor.sock ; (default is to run a UNIX domain socket server)

logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)

pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

user=nginx                   ; (default is current user, required if root)

[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[group:theq]
programs=theq-8001,theq-8002,theq-8003,theq-8004

[program:theq-8001]
command=/wherever/it/lives/app.py --port=8001
directory=/wherever/it/lives
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/theq-8001
logfile_maxbytes=1MB
logfile_backups=10

[program:theq-8002]
command=/wherever/it/lives/app.py --port=8002
directory=/wherever/it/lives
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/theq-8002
logfile_maxbytes=1MB
logfile_backups=10

[program:theq-8003]
command=/wherever/it/lives/app.py --port=8003
directory=/wherever/it/lives
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/theq-8003
logfile_maxbytes=1MB
logfile_backups=10

[program:theq-8004]
command=/wherever/it/lives/app.py --port=8004
directory=/wherever/it/lives
autostart=true
autorestart=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/theq-8004
logfile_maxbytes=1MB
logfile_backups=10
